# haspiaudio/rootfs/etc/services.d/minidlna/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Start minidlna' 
exec /usr/sbin/minidlnad -S  < /dev/null