# haspiaudio/rootfs/etc/services.d/snapserver/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Start SnapServer' 
exec  /usr/bin/snapserver < /dev/null

